# scriptplusconf/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command= /$REPONAME/scriptplusconf/nginx.sh
autorestart=true

[program:downscaled]
command= /$REPONAME/scriptplusconf/start.sh
autorestart=true
